# deploy/huggingface/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid
childlogdir=/tmp/supervisor
strip_ansi=true

[unix_http_server]
file=/tmp/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[program:redis]
command=/srv/deploy/service-entrypoint.sh redis /usr/bin/redis-server /srv/deploy/redis.conf
priority=10
autostart=true
autorestart=true
startsecs=2
startretries=20
stopsignal=TERM
stopwaitsecs=30
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:backend]
command=/srv/deploy/service-entrypoint.sh backend /opt/backend-venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 7860 --proxy-headers --forwarded-allow-ips=*
directory=/srv/backend
environment=PYTHONPATH="/srv/backend:/srv"
priority=20
autostart=true
autorestart=true
startsecs=3
startretries=20
stopsignal=TERM
stopwaitsecs=30
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:ai-worker]
command=/srv/deploy/service-entrypoint.sh ai-worker /opt/ai-venv/bin/celery -A app.worker:celery_app worker --loglevel=INFO --queues=ai --concurrency=1 --prefetch-multiplier=1 --hostname=ai@%%h
directory=/srv/ai_core
environment=PYTHONPATH="/srv/ai_core:/srv"
priority=30
autostart=true
autorestart=true
startsecs=5
startretries=20
stopsignal=TERM
stopwaitsecs=120
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:render-worker]
command=/srv/deploy/service-entrypoint.sh render-worker /opt/ai-venv/bin/celery -A app.worker:celery_app worker --loglevel=INFO --queues=render --concurrency=1 --prefetch-multiplier=1 --hostname=render@%%h
directory=/srv/ai_core
environment=PYTHONPATH="/srv/ai_core:/srv"
priority=40
autostart=true
autorestart=true
startsecs=5
startretries=20
stopsignal=TERM
stopwaitsecs=120
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
